# old/systemd/publish-ip.timer
[Unit]
Description=Publish dynamic IP

[Timer]
OnCalendar=*-*-* *:00/5:05
Persistent=true

[Install]
WantedBy=timers.target
