# systemd/lila-ws.service
[Unit]
Description=lila-ws
Requires=bloop.service

[Service]
ExecStart=bloop run root -m lila.ws.Boot -c /home/thib/lila-ws/.bloop

[Install]
WantedBy=multi-user.target
