# systemd/lila-ws.service
[Unit]
Description=lila-ws

[Service]
WorkingDirectory=/home/thib/lila-ws
ExecStart=/home/thib/.local/share/coursier/bin/sbtn reStart
ExecStop=/home/thib/.local/share/coursier/bin/sbtn reStop

[Install]
WantedBy=multi-user.target
