# systemd/lila-ws.service
[Unit]
Description=lila-ws
After=network.target

[Service]
Environment="JAVA_OPTS=-Xms128m -Xmx256m"
ExecStart=/home/thib/lila-ws/target/universal/stage/bin/lila-ws -Duse.epoll=true

[Install]
WantedBy=multi-user.target
